# extra/service/legion-linux.path
[Path]
Unit=legion-linux.service
PathModified=/sys/firmware/acpi/platform_profile

[Install]
WantedBy=default.target
