# extra/service/legiond.service
[Unit]
Description=legiond daemon

[Service]
ExecStart=/usr/bin/legiond

[Install]
WantedBy=multi-user.target
