# extra/service/legion-linux-restart.service
[Service]
Type=oneshot
ExecStart=/bin/systemctl restart legion-linux.service