# extra/service/legiond-cpuset.service
[Unit]
Description=set legion cpu

[Service]
Type=oneshot
ExecStart=/usr/bin/legiond-cli cpuset
