# extra/service/legion-linux-restart.path
[Path]
Unit=legion-linux-restart.service
PathModified=/sys/firmware/acpi/platform_profile

[Install]
WantedBy=default.target
